# CRM/server/launchd/com.vernontm.crm-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vernontm.crm-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/raysmacbook/Desktop/Vernon Tech And Media/Client Projects/VTM/CRM/server/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/raysmacbook/Desktop/Vernon Tech And Media/Client Projects/VTM/CRM/server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/raysmacbook/Library/Logs/vtm-crm-server.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/raysmacbook/Library/Logs/vtm-crm-server.err.log</string>
</dict>
</plist>
